# CodeDeck/CodeDeckService/codedeck.service
[Unit]
Description=CodeDeck Neural Interface - Local LLM Inference Hub
Documentation=https://github.com/codemusic/codedeck
After=network.target
Wants=network.target

[Service]
Type=simple
User=codemusic
Group=codemusic
WorkingDirectory=/home/codemusic/CodeDeck/CodeDeckService
Environment=PATH=/home/codemusic/CodeDeck/codedeck_venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/codemusic/CodeDeck/CodeDeckService
ExecStart=/home/codemusic/CodeDeck/codedeck_venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/codemusic/CodeDeck

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codedeck

[Install]
WantedBy=multi-user.target 